# JAP_WS_REMOTE_NODE/JAP_WS_REMOTE_NODE_DOTCLOUD/supervisord.conf
[program:node]
command = node JAP_WS_REMOTE.js
directory = /home/dotcloud/current