# JAP_REMOTE_WS_NODE/JAP_REMOTE_WS_NODE_DOTCLOUD/supervisord.conf
[program:node]
command = node JAP_REMOTE_WS.js
directory = /home/dotcloud/current